# projects/camstreamer/scripts/100-actioncamstreamer.rules
ENV{DEVTYPE}=="disk", ENV{SUBSYSTEM}=="block", ACTION=="add", RUN+="/actioncamstreamer/udev_start_usb_install.sh"