# projects/camstreamer/scripts/100-actioncamstreamer.rules
ENV{DRIVERS}=="usb-storage", ENV{SUBSYSTEM}=="usb", ACTION=="add", RUN+="/actioncamstreamer/udev_start_usb_install.sh"