# projects/camstreamer/scripts/start_actioncamstreamer.service
[Unit]
Description=ActionCamStreamer Service
After=multi-user.target

[Service]
Type=idle
ExecStart=/actioncamstreamer/start_camstreamer.sh

[Install]
WantedBy=multi-user.target
