# projects/camstreamer/scripts/start_actioncamstreamer.service
[Unit]
Description=ActionCamStreamer Service
Before=multi-user.target # Ensures that this unit is started before multi-user.target, otherwise they would start simultaneously. This use to be 'After', but the start time after a reboot was over a minute on the 2018 raspbian image.

[Service]
Type=idle
ExecStart=/.actioncamstreamer/start_actioncamstreamer.sh

# used by the enable/disable to link and unlink with dependencies
[Install]
WantedBy=multi-user.target
